# supervisor.conf
[program:check_status]
command=php /var/www/html/httpstatus/console.php checkStatus
autostart=true
autorestart=true